# vetscan-hub-platform-ubuntu/SetupScripts_v7.0/service_delayed_splash.sh
#!/bin/sh
# chkconfig: 345 99 10
case "$1" in
  start)
    # Executes our script
    sudo sh /home/vetscan/delayed_splash_exec.sh
    ;;
  *)
    ;;
esac
exit 0
